# swhkd/swhkd.service
[Unit]
Description=swhkd hotkey daemon
BindsTo=default.target

[Service]
Type=simple
ExecStart=/usr/local/lib/swhkd.sh

[Install]
WantedBy=default.target
